# systemd/lianli-daemon.service
[Unit]
Description=Lian Li Device Daemon
After=default.target

[Service]
Type=simple
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env lianli-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
